# cache/fetch_ids/src/supervisor/cluster_stuff.conf
[program:cluster_stuff]
environment=PYTHONPATH=/home/simon/.local/lib/python3.9/site-packages
command= /home/simon/Downloads/cdhit-4.8.1/cd-hit -i /home/simon/git/GroovIO/create_db/TetR/NCBI_TetRs_100-300.fasta -o /home/simon/git/GroovIO/create_db/TetR/50perc_clusters_trimmed -c 0.5 -n 3 -M 3000 -T 4
directory=/home/simon/git/GroovIO/create_db/TetR
autostart=true
autorestart=true
startretries=1
stderr_logfile=/var/log/opSeq/cluster_stuff.err.log
stdout_logfile=/var/log/opSeq/cluster_stuff.out.log
user=root
