# cache/fetch_ids/src/supervisor/create_db.conf
[program:create_db]
environment=PYTHONPATH=/home/simon/.local/lib/python3.9/site-packages
command= /usr/bin/python3 /home/simon/git/GroovIO/create_db/query2db.py
directory=/home/simon/git/GroovIO/create_db
autostart=true
autorestart=true
startretries=1
stderr_logfile=/var/log/opSeq/create_db.err.log
stdout_logfile=/var/log/opSeq/create_db.out.log
user=root
