# cache/fetch_ids/src/supervisor/harvest_operators.conf
[program:harvest_operators]
environment=BLASTPATH=/home/simon/ncbi-blast-2.13.0+/bin,PYTHONPATH=/home/simon/.local/lib/python3.9/site-packages
command= /usr/bin/python3 /home/simon/git/GroovIO/main.py
directory=/home/simon/git/GroovIO/
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/opSeq/harvest_operators.err.log
stdout_logfile=/var/log/opSeq/harvest_operators.out.log
user=root
